# scripts/com.housing-investment.db-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.housing-investment.db-backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/marioortizcuadrado/Desktop/Projects/housing-investment-api/scripts/backup-db.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/marioortizcuadrado/Desktop/Projects/housing-investment-api/backups/backup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/marioortizcuadrado/Desktop/Projects/housing-investment-api/backups/backup.log</string>
</dict>
</plist>
